# files/codeserver.conf
[program:codeserver]
command=/usr/bin/code-server
autorestart=true
autostart=true