# files/ngrok-code.conf
[program:ngrok-code]
command=/snap/bin/ngrok http 8080
autorestart=true
autostart=true