# files/schedulers.conf
[program:schedulers]
command=/usr/bin/python3 /home/ubuntu/Projects/Schedulers/pi_schedule.py
autorestart=true
autostart=true